# transcribe-cpp/rootfs/etc/s6-overlay/s6-rc.d/transcribe-cpp/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Service finish handler: every service exit halts the container.
#
# Clean exits and signals (manual stop) halt with their code.  Bootstrap /
# configuration failures (bad custom_model, failed conversion, ...) log the
# reason and halt with the service's exit code — the add-on stops instead of
# lingering half-alive or silently retrying; logs stay readable in the HA UI.
# ==============================================================================
# shellcheck disable=SC2155
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"

bashio::log.info \
  "Service exited with code ${exit_code_service}" \
  "(by signal ${exit_code_signal})"

if [[ "${exit_code_signal}" -ne 0 ]]; then
  # Killed by signal — transient (manual stop / resource limit).  Halt.
  if [[ "${exit_code_container}" -eq 0 ]]; then
    echo $((128 + exit_code_signal)) > /run/s6-linux-init-container-results/exitcode
  fi
  exec /run/s6/basedir/bin/halt
elif [[ "${exit_code_service}" -eq 0 ]]; then
  # Clean exit — transient (graceful shutdown).  Halt.
  exec /run/s6/basedir/bin/halt
else
  # Non-zero exit — configuration/bootstrap failure (bad custom_model,
  # failed conversion, model download, ...).  Stop the whole add-on.
  bashio::log.error \
    "transcribe-cpp: fatal error (exit ${exit_code_service}) —" \
    "stopping the add-on. Fix the configuration and start it again."
  if [[ "${exit_code_container}" -eq 0 ]]; then
    echo "${exit_code_service}" > /run/s6-linux-init-container-results/exitcode
  fi
  exec /run/s6/basedir/bin/halt
fi
